# flask-status.service
[Unit]
Description=Flask Status Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/pi-scripts/flask_status.py
Restart=on-failure
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
