# wifi-ap.service
[Unit]
Description=Wi-Fi Access Point on wlan1
After=network-pre.target
Before=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/pi-scripts/setup_ap.py
Restart=on-failure
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
